# systemd/findex-scanner.timer
[Unit]
Description=Findex Scanner Timer - Periodic File Index Updates
Documentation=https://github.com/ogefest/findex

[Timer]
# Run 15 minutes after boot
OnBootSec=15min

# Run every hour
OnCalendar=hourly

# Small random delay to avoid exact hour load spikes
RandomizedDelaySec=300

# Persist timer across reboots
Persistent=true

[Install]
WantedBy=timers.target
